# raspberry/shutdown-button.service
[Unit]
Description=Shutdown when button is pressed

[Service]
ExecStart=/home/alarm/on-shutdown.py

[Install]
WantedBy=multi-user.target

